# btrfs-lightning-search/btrfs-indexer.service
[Unit]
Description=Btrfs File Indexer inotify Daemon
Documentation=https://github.com/user/btrfs_file_searcher
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/local/share/btrfs-lightning-search
ExecStart=/usr/bin/python3 /usr/local/share/btrfs-lightning-search/inotify_daemon.py
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/USER_HOME/.local/share/btrfs-search

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btrfs-indexer

[Install]
WantedBy=multi-user.target
